# control/units/logspout.service
[Unit]
Description=logspout
Requires=docker.service
After=docker.service

[Service]
Type=simple
Restart=on-failure
RestartSec=5

EnvironmentFile=/etc/environment
Environment="ENV=blue"
Environment="TAG=latest"

TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker kill logspout
ExecStartPre=-/usr/bin/docker rm logspout
ExecStartPre=/usr/bin/docker pull progrium/logspout:${TAG}

ExecStart=/usr/bin/docker run \
  --name logspout \
  --hostname %m \
  --dns-search ${ENV}.cluster.local \
  --env DEBUG=true \
  --volume /var/run/docker.sock:/tmp/docker.sock \
  --publish 8000:8000 \
  progrium/logspout:${TAG} \
    syslog://172.17.42.1:514


[X-Fleet]
Global=true
MachineMetadata="role=control"

[Install]
WantedBy=multi-user.target
